# scripts/systemd/katana-api.service
[Unit]
Description=Katana API - .NET 8 Backend Service
After=network.target
Wants=network-online.target
Documentation=https://github.com/gamzeaydinnn/katana

[Service]
Type=simple
User=huseyinadm
Group=huseyinadm
WorkingDirectory=/home/huseyinadm/katana/publish
ExecStart=/usr/bin/dotnet /home/huseyinadm/katana/publish/Katana.API.dll
Restart=on-failure
RestartSec=10
TimeoutStopSec=90
SyslogIdentifier=katana-api

# Environment
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false
Environment=ASPNETCORE_URLS=http://0.0.0.0:5055

# Security hardening (optional but recommended)
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
